# apps/admin-web/scripts/local-cron/plists/com.nepalrepublic.cron-sweep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Full intelligence sweep — runs all collectors (RSS / YouTube / social), classifies via local AI.

  Schedule: 07:15 and 19:15 NPT (twice a day, matching prior Vercel cron times).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nepalrepublic.cron-sweep</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-c</string>
    <string>"$HOME/Library/Application Support/nepal-republic/cron-runner.sh" sweep-full '/api/intelligence/sweep?mode=full'</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>15</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>19</integer>
      <key>Minute</key><integer>15</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/priyanka.shrestha/Library/Logs/nepal-republic/cron-sweep.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/priyanka.shrestha/Library/Logs/nepal-republic/cron-sweep.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/priyanka.shrestha</string>
    <key>PATH</key>
    <string>/Users/priyanka.shrestha/.nvm/versions/node/v20.20.0/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
